# deploy/claude-bot.service
# OCI 서버(Rocky Linux)에서 봇을 상시 실행하는 systemd 서비스
# 설치: sudo cp deploy/claude-bot.service /etc/systemd/system/
#       sudo systemctl daemon-reload && sudo systemctl enable --now claude-bot
#
# SELinux 주의 (Rocky/RHEL):
#  - EnvironmentFile은 홈 디렉터리에 두면 거부됨 → /etc/claude-bot.env 사용
#    (파일을 mv로 옮겼다면 `sudo restorecon /etc/claude-bot.env` 로 라벨 복구)
#  - 홈 디렉터리의 .venv/bin/python 직접 실행도 거부됨 → 시스템 python + PYTHONPATH 사용
[Unit]
Description=Telegram Claude Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rocky
WorkingDirectory=/home/rocky/telegram-bot
EnvironmentFile=/etc/claude-bot.env
Environment=PYTHONPATH=/home/rocky/telegram-bot/.venv/lib/python3.12/site-packages
ExecStart=/usr/bin/python3.12 /home/rocky/telegram-bot/bot.py
Restart=always
RestartSec=10

# balruno 격리: 봇(코딩 모드)이 balruno 데이터/설정을 절대 못 건드리게 차단.
# 마운트 네임스페이스라 sudo로도 우회 불가. balruno 서비스 자체엔 영향 없음.
InaccessiblePaths=/opt/balruno /var/lib/balruno /var/lib/balruno-backup /var/log/balruno

[Install]
WantedBy=multi-user.target
